# collector/gpu-metrics-collector.service
[Unit]
Description=GPU Metrics Collector
Documentation=https://github.com/your-org/gpu-metrics-collector
After=network.target
Wants=network.target

[Service]
Type=simple
User=gpu-metrics
Group=gpu-metrics
ExecStart=/usr/local/bin/gpu-metrics-collector
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/gpu-metrics
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Environment
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory=/var/lib/gpu-metrics

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gpu-metrics-collector

# Resource limits
LimitNOFILE=1024
LimitNPROC=512

[Install]
WantedBy=multi-user.target
